# etc/inittab
# /etc/inittab: init configuration for busybox init.
# Boot-time system configuration/initialization script.
#
::sysinit:/etc/init.d/rcS
ttymxc0::respawn:/bin/sh
tty1::respawn:/bin/sh
